# ops/pulsar-sqlite-backup.timer
[Unit]
Description=Run Pulsar SQLite backup every day

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=pulsar-sqlite-backup.service

[Install]
WantedBy=timers.target
